# com.eshanb.tesla-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eshanb.tesla-tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/eshanb/workspace/projects/tesla-tracker/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/eshanb/workspace/projects/tesla-tracker</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/eshanb/workspace/projects/tesla-tracker/tracker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eshanb/workspace/projects/tesla-tracker/tracker.log</string>
</dict>
</plist>
